# quadletman.service
[Unit]
Description=quadletman - Podman Quadlet web manager
Documentation=https://github.com/mikkovihonen/quadletman
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/quadletman
WorkingDirectory=/var/lib/quadletman
StateDirectory=quadletman
LogsDirectory=quadletman
RuntimeDirectory=quadletman
Environment=QUADLETMAN_LOG_LEVEL=INFO
Restart=on-failure
RestartSec=5
# Security hardening (limited — root required for user management)
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
